# public_ip_monitoring.timer
[Unit]
Description=Schedule public IP monitoring hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

